# etc/rc.d/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local

# Need to add the following section

for each in /proc/sys/net/ipv4/conf/*; do
    echo 0 > $each/accept_redirects
    echo 0 > $each/send_redirects
    echo 0 > $each/rp_filter
done

# May need the following but not sure

# echo 1 >  /proc/sys/net/core/xfrm_larval_drop